# boards/vsrves01/rootfs-overlay/etc/inittab
# /etc/inittab for VSRVES01 (busybox/OpenRC).
#
# Bootargs from catboard chosen.bootargs put kernel console on ttyUL0@1000000
# (UART-Lite emulated). The 'term -a' wrapper around the simulator forwards
# stdio onto the same line, so we also enable a 115200 ttyS0 getty for
# environments where stdio is wired through a more conventional UART.

# system init -> openrc
::sysinit:/sbin/openrc sysinit
::sysinit:/sbin/openrc boot
::wait:/sbin/openrc default

# UART-Lite console (matches kernel cmdline)
ttyUL0::respawn:/sbin/getty -L 1000000 ttyUL0 vt100

# Conventional serial fallback
ttyS0::respawn:/sbin/getty -L 115200 ttyS0 vt100

# Stuff to do before rebooting
::shutdown:/sbin/openrc shutdown

# Ctrl-Alt-Del disabled (no PC keyboard on this board anyway)
#::ctrlaltdel:/sbin/reboot
